# files/profile.d/ls-aliases.sh
alias ls="ls --color --classify"
alias ll="ls --color --classify --human-readable -l"
