# files/profile.d/path-setup.sh
export PATH="${PATH}:/usr/bin/core_perl"

# Use ccache's compilers to get caching
export PATH="${PATH}:/usr/lib/ccache/bin"
