# files/profile.d/sway.sh
WLC_REPEAT_DELAY=180
WLC_REPEAT_RATE=35
