# files/profile.d/systemctl-aliases.sh
alias sysc="systemctl"
alias jrnc="journalctl"
alias sysu="systemctl --user"
alias jrnu="journalctl --user"
