# templates/profile.d/base16-theme.sh.j2
export MY_BASE16_THEME={{ base16_theme }}
