# files/libinput-gestures.service
[Unit]
Description=Multitouch via libinput-gestures

[Service]
ExecStart=/usr/bin/libinput-gestures

[Install]
WantedBy=default.target
