# files/libinput-gestures.service
[Unit]
Description=Multitouch via libinput-gestures
BindsTo=xorg.service

[Service]
Environment=DISPLAY=:0 # TODO: pull this out
ExecStart=/usr/bin/libinput-gestures

[Install]
WantedBy=xorg.service
