# files/package-update.service
[Unit]
Description=Update installed packages

[Service]
Type=oneshot
ExecStart=/usr/sbin/pacaur -Syu

[Install]
RequiredBy=default.target
