# files/pass-credential@.service
[Unit]
Description=pass credential server for %i
Requires=pass-credential@%i.socket
After=pass-credential@%i.socket
StartLimitIntervalSec=0

[Service]
ExecStart=/usr/bin/pass-systemd-cred-server --credential %i

# The user session environment (PASSWORD_STORE_DIR, GNUPGHOME, gpg-agent
# socket) is inherited automatically. Uncomment to scope all lookups under
# a sub-directory of the store:
# Environment=PASS_CREDENTIAL_PREFIX=systemd

NoNewPrivileges=yes
# full (not strict) keeps the home directory writable, which gpg may need
# for trust-db updates and lock files.
ProtectSystem=full
PrivateTmp=yes
