# files/pass-credential@.socket
[Unit]
Description=pass credential socket for %i

[Socket]
# %t expands to $XDG_RUNTIME_DIR (e.g. /run/user/1000).
# %i is the credential name — the part between @ and .socket in the unit name.
# Socket is placed directly in %t to avoid RuntimeDirectory sharing conflicts
# across multiple template instances.
ListenStream=%t/pass-credential-%i.socket
SocketMode=0600

[Install]
WantedBy=default.target
