# files/sway-session.target
[Unit]
Description=Sway session

After=graphical-session-pre.target
BindsTo=graphical-session.target
Wants=graphical-session-pre.target
