# files/x11.socket
[Unit]
Description=Socket for xorg

[Socket]
ListenStream=/tmp/.X11-unix/X
Service=xorg.service

[Install]
WantedBy=sockets.target
