# files/x11.target
[Unit]
Description=X11 target

After=xrdb.service
Requires=xrdb.service

[Install]
WantedBy=default.target
