# files/x11@.target
[Unit]
Description=X11 target for display %i

After=xrdb@%i.service
Requires=xrdb@%i.service

[Install]
WantedBy=default.target
